# deploy/systemd/tg-kb-alert@.service
[Unit]
Description=Alert hook for failed unit %i

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'logger -p user.err -t tg-kb-alert "unit failed: %i (see: systemctl status %i && journalctl -u %i -n 200)"'
