# deploy/systemd/tg-kb-healthcheck.service
[Unit]
Description=Runtime healthcheck for tg-kb-bot

[Service]
Type=oneshot
User=tgkb
Group=tgkb
WorkingDirectory=/opt/tg_db
ExecStart=/usr/local/bin/tg-kb-runtime-healthcheck
